# supervisor.conf
[program:bloodhunt]
command = /usr/local/sbin/uwsgi --ini /www/bloodhunt/uwsgi.conf
process_name=uwsgi
directory = /www/bloodhunt/src
user=www

